# templates/default/node_exporter.service.erb
[Unit]
Description = Node exporter 
After = syslog.target network.target remote-fs.target
StartLimitIntervalSec=60
StartLimitBurst=20

[Service]
User = <%= node['hopsmonitor']['user'] %> 
Group = <%= node['hopsmonitor']['group'] %> 

ExecStart = <%= node['node_exporter']['base_dir'] %>/node_exporter \
            --web.listen-address=:<%= node['node_exporter']['port'] %> \
            --collector.filesystem.ignored-mount-points=<%= node['node_exporter']['filesystem']['regex'] %> \
            --collector.textfile.directory=<%= node['node_exporter']['text_metrics'] %>

LimitNOFILE = 20000
Restart=on-failure
RestartSec=2

Type = simple

[Install]
WantedBy = multi-user.target