# templates/default/prometheus.service.erb
[Unit]
Description = Prometheus  
After = syslog.target network.target remote-fs.target
StartLimitIntervalSec=60
StartLimitBurst=20

[Service]
User = <%= node['hopsmonitor']['user'] %> 
Group = <%= node['hopsmonitor']['group'] %> 

ExecStart = <%= node['prometheus']['base_dir'] %>/prometheus \
            --config.file=<%= node['prometheus']['base_dir'] %>/prometheus.yml \
            --web.listen-address=0.0.0.0:<%= node['prometheus']['port'] %> \
            --storage.tsdb.path=<%= node['prometheus']['data_dir'] %> \
            --storage.tsdb.retention.time=<%= node['prometheus']['retention_time'] %>
LimitNOFILE = 20000
Restart=on-failure
RestartSec=2

Type = simple

[Install]
WantedBy = multi-user.target